# docker/crontab
*/5 * * * * /usr/local/bin/cloudflare-ddns
