# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS

[Service]
Type=oneshot
EnvironmentFile=-/etc/cloudflare-ddns
ExecStart=/usr/local/bin/cloudflare-ddns
